# deploy/django.service
[Unit]
Description=AI Web Ollama - Django (Gunicorn)
After=network.target postgresql.service

[Service]
User=admin
WorkingDirectory=/home/admin/ai_web_ollama/backend
EnvironmentFile=/home/admin/ai_web_ollama/backend/.env
ExecStart=/home/admin/ai_web_ollama/backend/venv/bin/gunicorn \
    config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 2 \
    --timeout 300 \
    --access-logfile /home/admin/ai_web_ollama/logs/access.log \
    --error-logfile /home/admin/ai_web_ollama/logs/error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
